# templates/default/statsd.service.erb
[Unit]
Description=Statsd Librato

[Service]
User=<%= node[:statsd][:user] %>
Environment=LIBRATO_USER=<%= node[:statsd][:librato_user] %> 
Environment=LIBRATO_TOKEN=<%= node[:statsd][:librato_token] %> 
Enviornment=LIBRATO_SOURCE=<%= node[:statsd][:librato_source] %>

ExecStart=/usr/local/bin/statsd <%= @args.map {|k,v| "-#{k}=#{v}" }.join(" ") %>
Restart=on-abort

[Install]
WantedBy=default.target
